# ansible/roles/package-daemons/files/metric-reporting-daemon.service
[Unit]
Description=Metric Reporting Daemon - Provides HTTP API to allow injecting custom metrics into the metric stream.

[Service]
Environment=JWT_SECRET=8a43b51b2d5ea077c5c3b8e1f84d1f967630bb548a0eb5800e2805b4a0c9beed7d5245bf408dd33954d36691b5612c3fa08b00fadd1395e722c67c445a59d98ab76c7fdc573b4f91c41e5e213c482e500f02e6b8d36e6ef2666959da2a1201cdcd256aae6ece101533b98f73a2209fdf55047d68b5b897f793035da7694a666174d773fe619fb645b557ac3891d51421f543a45d47112615db5ce01a190aaf6f3fb605f6da28b274be8c995f5d26be7b6b8de7b17d3fcb1f30c2c73da4c5f5052004ea069ae1d3721fa4f2838108a3b5b0fa0501610d3f52de9bc080a00cde7930d3ad5a9bcee2f16636a50532c291b29f21aef0e71957b52a877b32f2ace810
Type=exec
ExecStart=/data/private/share/daemons/metric-reporting-daemon/metric-reporting-daemon
TimeoutSec=10sec
StandardOutput=append:/var/log/metric-reporting-daemon.log
StandardError=append:/var/log/metric-reporting-daemon.log

[Install]
WantedBy=multi-user.target
